# PycharmProjects/web13/webSecond/bbs.conf
[program:bbs]
command=/root/webSecond/gunicorn wsgi -c gunicorn.config.py
directory=/root/webSecond
autostart=true
autorestart=true
